# systemd-timers/daily-report.service
[Unit]
Description=Generate Daily Reward Report
After=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/var/www/daily-rewards
EnvironmentFile=/var/www/daily-rewards/.env
ExecStart=/usr/bin/pm2 restart daily-report --stop-exit-codes 0
StandardOutput=append:/var/log/daily-report.log
StandardError=append:/var/log/daily-report.log
TimeoutStartSec=300

# Retry configuration
Restart=on-failure
RestartSec=30
StartLimitBurst=5
StartLimitIntervalSec=600

[Install]
WantedBy=multi-user.target
