# systemd-timers/entry-task.service
[Unit]
Description=Daily Rewards Entry Task
After=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/var/www/daily-rewards
EnvironmentFile=/var/www/daily-rewards/.env

# Run entry task via the new orchestrator CLI
# Set FEATURE_MODE=production for all features, or development/dry-run for testing
ExecStart=/usr/bin/node node_modules/.bin/ts-node ./src/modules/orchestrator/cli.ts

StandardOutput=append:/var/log/entry-task.log
StandardError=append:/var/log/entry-task.log
TimeoutStartSec=7200

# Retry configuration
Restart=on-failure
RestartSec=60
StartLimitBurst=3
StartLimitIntervalSec=600

[Install]
WantedBy=multi-user.target
