# systemd-timers/entry-task.service
[Unit]
Description=Daily Rewards Entry Task
After=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/var/www/daily-rewards
EnvironmentFile=/var/www/daily-rewards/.env
ExecStart=/usr/bin/node node_modules/.bin/ts-node ./src/modules/entry.ts
StandardOutput=append:/var/log/entry-task.log
StandardError=append:/var/log/entry-task.log
TimeoutStartSec=300

# Retry configuration
Restart=on-failure
RestartSec=30
StartLimitBurst=5
StartLimitIntervalSec=600
