# systemd-timers/entry-task.timer
[Unit]
Description=Daily Rewards Entry Task Timer

[Timer]
OnCalendar=*-*-* 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
