# systemd-timers/entry-task.service
[Unit]
Description=Daily Rewards Entry Task
After=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/var/www/daily-rewards
EnvironmentFile=/var/www/daily-rewards/.env

Environment=FEATURE_MODE=production
ExecStart=/usr/bin/flock -n /var/lock/daily-rewards.lock /usr/bin/node dist/modules/orchestrator/cli.js
StandardOutput=append:/var/log/entry-task.log
StandardError=append:/var/log/entry-task.log
TimeoutStartSec=7200

# Retry configuration
Restart=on-failure
RestartSec=60
StartLimitBurst=3
StartLimitIntervalSec=600

[Install]
WantedBy=multi-user.target
