# systemd-timers/unpause-task.service
[Unit]
Description=Daily Rewards Unpause Task
After=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/var/www/daily-rewards
EnvironmentFile=/var/www/daily-rewards/.env

# Run unpause task directly via ts-node
ExecStart=/usr/bin/node node_modules/.bin/ts-node ./src/modules/unpause.ts

StandardOutput=append:/var/log/unpause-task.log
StandardError=append:/var/log/unpause-task.log
TimeoutStartSec=300

# Retry configuration
Restart=on-failure
RestartSec=30
StartLimitBurst=5
StartLimitIntervalSec=600

[Install]
WantedBy=multi-user.target
