# systemd-timers/unpause-task.service
[Unit]
Description=Daily Rewards Unpause Task
After=network.target entry-task.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/var/www/daily-rewards
EnvironmentFile=/var/www/daily-rewards/.env

# Run guarded unpause task directly from the built output
ExecStart=/usr/bin/flock -n /var/lock/daily-rewards.lock /usr/bin/node dist/modules/unpause.js

StandardOutput=append:/var/log/unpause-task.log
StandardError=append:/var/log/unpause-task.log
TimeoutStartSec=300

# Retry configuration
Restart=on-failure
RestartSec=30
StartLimitBurst=5
StartLimitIntervalSec=600

[Install]
WantedBy=multi-user.target
