# systemd-timers/unpause-task.timer
[Unit]
Description=Daily Rewards Unpause Task


[Timer]
OnCalendar=*:40
Persistent=true

[Install]
WantedBy=timers.target
