# systemd-timers/unpause-task.timer
[Unit]
Description=Daily Rewards Unpause Task Timer

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
